# src/etc/profile.d/stardog
STARDOG_HOME=/opt/stardog
STARDOG=/var/lib/stardog
